# systemd/user/flatpak-update.timer
[Unit]
Description=Timer per aggiornare i Flatpak quotidianamente

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
